# com.goofish.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.goofish.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/weiyi/workspaces/ai-goofish-monitor/venv/bin/python3</string>
        <string>web_server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/weiyi/workspaces/ai-goofish-monitor</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/weiyi/workspaces/ai-goofish-monitor/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/weiyi/workspaces/ai-goofish-monitor/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
